# deploy/cloud-clipboard.service
[Unit]
Description=Cloud Clipboard - Image & Text Storage
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cloud-clipboard
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000

# Security hardening
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloud-clipboard

[Install]
WantedBy=multi-user.target
